# docker/supervisord.conf
[supervisord]
nodaemon=true
# /tmp is world-writable, so supervisord can write its pidfile while running as
# the image's non-root www-data user (no need to chown /var/run).
pidfile=/tmp/supervisord.pid
logfile=/dev/fd/1
logfile_maxbytes=0
loglevel=info

# Octane web server (FrankenPHP worker mode) — the primary process.
# stop/killasgroup ensure SIGTERM reaches the FrankenPHP worker children on
# `docker stop` / supervisord shutdown (graceful Octane drain).
[program:octane]
command=php /var/www/html/artisan octane:start --server=frankenphp --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# Queue worker (database connection by default). TERM lets queue:work finish the
# in-flight job before exiting; stopwaitsecs gives it room to do so.
# Guarded: idles when QUEUE_WORKER_ENABLED=false — set that to run the worker as a
# separate service/container instead (e.g. in a cloud deployment).
[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=/bin/bash -c 'if [ "$QUEUE_WORKER_ENABLED" != "false" ]; then exec php /var/www/html/artisan queue:work --tries=3 --max-time=3600; else echo "[worker] disabled (QUEUE_WORKER_ENABLED=false); idling"; exec sleep infinity; fi'
autostart=true
autorestart=true
numprocs=%(ENV_QUEUE_WORKER_COUNT)s
stopsignal=TERM
stopwaitsecs=35
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# Scheduler. schedule:work forks due commands as child processes, so stop the
# whole group to avoid orphaning a running scheduled command.
# Guarded: idles when SCHEDULER_ENABLED=false — set that to run the scheduler as a
# separate service/container instead (e.g. in a cloud deployment).
[program:schedule]
command=/bin/bash -c 'if [ "$SCHEDULER_ENABLED" != "false" ]; then exec php /var/www/html/artisan schedule:work; else echo "[schedule] disabled (SCHEDULER_ENABLED=false); idling"; exec sleep infinity; fi'
autostart=true
autorestart=true
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# Inertia SSR — guarded: idles unless INERTIA_SSR_ENABLED=true (no flapping).
[program:ssr]
command=/bin/bash -c 'if [ "$INERTIA_SSR_ENABLED" = "true" ]; then exec php /var/www/html/artisan inertia:start-ssr --runtime=bun; else echo "[ssr] disabled (INERTIA_SSR_ENABLED != true); idling"; exec sleep infinity; fi'
autostart=true
autorestart=true
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
